# systemd/ogarden.timer
[Unit]
Description=Ogarden Python Application
Documentation=https://github.com/drewlio/ogarden

[Timer]
# Run the task at these times
OnCalendar=06:00,18:00
Persistent=true

[Install]
WantedBy=timers.target
